# apps/flutter/ghostty/dist/linux/systemd.service.in
[Unit]
Description=@NAME@
After=graphical-session.target
After=dbus.socket
Requires=dbus.socket

[Service]
Type=notify-reload
ReloadSignal=SIGUSR2
BusName=@APPID@
ExecStart=@GHOSTTY@ --gtk-single-instance=true --initial-window=false

[Install]
WantedBy=graphical-session.target
